# web/docker/supervisord.conf
[supervisord]
nodaemon=true

[program:nginx]
command=nginx -g "daemon off;"
autostart=true
autorestart=true
redirect_stderr=true

[program:flask]
directory=/app
command=gunicorn --bind 0.0.0.0:5000 tps-proxy-server:app
autostart=true
autorestart=true
redirect_stderr=true
